# data/autostart/org.vaxp.vstore.desktop.in
[Desktop Entry]
Name=Vstore
Comment=Software Updates
Icon=org.vaxp.vstore
OnlyShowIn=GNOME;Unity;
NotShowIn=Budgie
NoDisplay=true
Exec=vstore
